# ops/systemd/gotra-stock-pool-evening-report.service
[Unit]
Description=GOTRA public stock-pool evening/HK report
After=network-online.target gotra-public-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/gotra
ExecStart=/root/.local/bin/uv run python scripts/public_stock_pool_report.py --mode evening-hk --publish-static --allowed-missing-symbol HKEX:0501
Nice=5
