# ops/systemd/gotra-stock-pool-us-morning-report.service
[Unit]
Description=GOTRA public stock-pool US morning report
After=network-online.target gotra-public-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/gotra
ExecStart=/root/.local/bin/uv run python scripts/public_stock_pool_report.py --mode morning-us --publish-static
Nice=5
