# ops/systemd/gotra-stock-pool-evening-report.timer
[Unit]
Description=Run GOTRA public stock-pool evening/HK report after HK market close

[Timer]
OnCalendar=Mon..Fri *-*-* 18:30:00 Asia/Shanghai
Persistent=true
Unit=gotra-stock-pool-evening-report.service

[Install]
WantedBy=timers.target
